# mqtt-bridge/mqtt-bridge.service
[Unit]
Description=MQTT Bridge — forward solar MQTT data to InfluxDB via HTTP
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/mqtt-bridge.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true

[Install]
WantedBy=multi-user.target
